# ops/ise547-backend.conf
[program:ise547-backend]
directory=/home/yuyaoli/personal/ISE547/backend
command=/home/yuyaoli/.local/bin/uv run uvicorn app.main:app --host 0.0.0.0 --port 8024
user=yuyaoli
autostart=true
autorestart=true
startsecs=5
stopsignal=TERM
stopwaitsecs=30
stdout_logfile=/home/yuyaoli/personal/ISE547/backend/logs/backend.stdout.log
stderr_logfile=/home/yuyaoli/personal/ISE547/backend/logs/backend.stderr.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=5
environment=PYTHONUNBUFFERED="1"
